# deploy/linux/uclaw.service
[Unit]
Description=UCLAW Web Server
After=network.target

[Service]
Type=simple
User=uclaw
Group=uclaw
WorkingDirectory=/opt/uclaw
EnvironmentFile=/etc/uclaw/uclaw.env
ExecStartPre=/usr/bin/npm --prefix /opt/uclaw run deploy:check
ExecStart=/usr/bin/npm --prefix /opt/uclaw start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/uclaw/.uclaw
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
